# ICICI_UI/supervisor.conf
[supervisord]
nodaemon=true

[program:npm]
directory=/ICICI_UI
command=npm start
autostart = true
autorestart = false
